# crontabs
david@keyweb:~$ crontab -l
# m h  dom mon dow   command
MAILTO=david@cantrell.org.uk
0 9 * * * nice -n 19 /web/cpandeps.cantrell.org.uk/getfiles.sh
0 7 * * * find /web/cpandeps.cantrell.org.ukdb/META.yml/ -mtime 7 -exec rm -f {} \;
0 7 * * * find /web/cpandeps.cantrell.org.ukdb/MANIFEST/ -mtime 7 -exec rm -f {} \;
0 7 * * * find /web/cpandeps.cantrell.org.ukdb/reverse/ -mtime 7 -exec rm -f {} \;
0 7 * * * find /web/cpandeps.cantrell.org.ukdb/results/ -mtime 7 -exec rm -f {} \;
